# units/mpris-scrobbler.service
[Service]
Type = simple
ExecStart = %h/playground/mpris-scrobbler/mpris-scrobbler -v
ExecReload = /bin/kill -HUP $MAINPID

