# units/mpris-scrobbler.service.in
[Service]
Type = simple
ExecStart = DESTDIR`'INSTALL_PREFIX`'BINDIR mpris-scrobbler -v
ExecReload = /bin/kill -HUP $MAINPID

