# units/systemd-user.service.in
[Unit]
Description = Small daemon to submit music retrieved from MPRIS DBus interface to last.fm and compatible services
Requires = dbus.socket

[Service]
Type = simple
ExecStart = DESTDIR`'INSTALL_PREFIX`'BINDIR/mpris-scrobbler -v
ExecReload = /bin/kill -HUP $MAINPID

