# units/systemd-user.service.in
[Unit]
Description = daemon to scrobble tracks loaded from the MPRIS DBus interface to compatible services
Requires = dbus.socket

[Service]
Type = simple
ExecStart = BINPATH`'DAEMONNAME -vv
ExecReload = /bin/kill -HUP $MAINPID

